# ansible-setup.service
[Unit]
Description=Install Ansible Tower
After=local-fs.target

[Service]
ExecStart=/docker-entrypoint.sh

[Install]
WantedBy=default.target
